# henet.service
[Unit]
Description=Henet Daily Test
OnFailure=status-email-user@%n.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/henet.py
